# debian/logilab-packaging.cron.weekly
#!/bin/sh

/usr/bin/lgp setup -c update -d all -a any
