# debian/logilab-lgp.cron.daily
#!/bin/sh

if ! [ -x /usr/bin/lgp ]; then
    exit 0
fi

# Update image
/usr/bin/lgp setup -qq -c update -d all -a any
