# debian/logilab-packaging.cron.daily
#!/bin/sh

# Clean caches
/usr/bin/lgp setup -c clean -d all -a any

# Update image
/usr/bin/lgp setup -c update -d all -a any
