# debian/logilab-packaging.cron.weekly
#!/bin/sh -e

lgp setup -c update -d all -a all
