# debian/logilab-packaging.cron.daily
#!/bin/sh

# Update image
/usr/bin/lgp setup -qq -c update -d all -a any
